# stRevizor-Server/misc/revizor.cron
*/6 * * * /path/to/f.wall/f.wall.revizor
